# crontab
# Run every minute for testing
0 6,20 * * * cd /app && /usr/bin/python3 /app/main.py >> /app/cron.log 2>&1

# Empty line required at end of cron file
